# src/lib/systemd/system/raspivid.service
[Unit]
Description=raspivid
After=network.target

[Service]
IgnoreSIGPIPE=false
ExecStart=/usr/bin/camera-stream

[Install]
WantedBy=default.target
