# src/crontab
*/5 * * * * sudo systemctl start raspivid.service
